# .dev/profile.d/aliases.sh
alias c='clear'
alias ll='ls -lAh'
alias cl='clear && ls -lAh'
alias ..='cd ../'
alias ...='cd ../../'
